# deploy/systemd/phantom-lab-health-watchdog.timer
[Unit]
Description=Run Phantom Lab Chat Docker health watchdog every minute

[Timer]
OnBootSec=2min
OnUnitActiveSec=1min
AccuracySec=15s
Unit=phantom-lab-health-watchdog.service

[Install]
WantedBy=timers.target
