# deploy/systemd/phantom-lab-health-watchdog.service
[Unit]
Description=Phantom Lab Chat Docker health watchdog
Wants=docker.service
After=docker.service

[Service]
Type=oneshot
User=deploy
Group=docker
WorkingDirectory=/opt/phantom-lab-chat
Environment=PHANTOM_LAB_ROOT=/opt/phantom-lab-chat
Environment="PHANTOM_LAB_WATCHDOG_SERVICES=api web"
ExecStart=/opt/phantom-lab-chat/deploy/scripts/health-watchdog.sh
